# 99-sdgrabber.rules
ACTION=="add", KERNEL=="sd[a-z][0-9]", SUBSYSTEMS=="usb", RUN+="/usr/local/bin/mount-sdcard.sh '%k'"