# 99-sdgrabber.rules
KERNEL=="sd[a-z][0-9]", SUBSYSTEMS=="usb", ACTION=="add", RUN+="/usr/local/bin/mount-sdcard.sh '%k'"