# 99-sdgrabber.rules
ACTION=="add", KERNEL=="sd[a-z][0-9]", SUBSYSTEMS=="usb", TAG+="systemd", ENV{SYSTEMD_WANTS}="sdgrabber-mount@%k.service"
ACTION=="remove", KERNEL=="sd[a-z][0-9]", SUBSYSTEMS=="usb", RUN+="/usr/local/bin/unmount-sdcard.sh"